# systemd/mw-refreshSFS.service
[Unit]
Description=Refresh StopForumSpam filter for ShinyWiki

[Service]
Type=oneshot
ExecStart=/opt/wiki/StopForumSpam/cron.sh
User=http
Group=http
